# private_dot_config/systemd/exact_user/cameractrlsd.service
[Unit]
Description=CameraCtrls daemon - restore control values

[Service]
Type=simple
ExecStart=/usr/bin/cameractrlsd

[Install]
WantedBy=graphical-session.target
