# private_dot_config/systemd/exact_user/rustup-update.timer
[Unit]
Description=Check for rust updates every day.

[Timer]
OnCalendar=*-*-* 13:00
Persistent=true

[Install]
WantedBy=timers.target
